# junctionbackend/deployment/supervisor.ini
[supervisord]
loglevel = DEBUG
logfile = /tmp/supervisord.log

[program:uwsgi-app]
redirect_stderr=true
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
command = /usr/local/bin/uwsgi --ini /app/deployment/uwsgi.ini

[program:nginx-app]
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
command = /usr/sbin/nginx -g "daemon off;"

